# systemd/wlan0-ap.service
[Unit]
Description=Configure wlan0 for Access Point
Before=hostapd.service
After=NetworkManager.service network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Configurar interfaz wlan0 (NetworkManager ya excluido en /etc/NetworkManager/conf.d/)
ExecStart=/sbin/ip addr flush dev wlan0
ExecStart=/sbin/ip addr add 192.168.50.1/24 dev wlan0
ExecStart=/sbin/ip link set wlan0 up

# Verificar que quedó UP
ExecStartPost=/bin/sh -c 'sleep 1 && ip link show wlan0 | grep -q "state UP"'

[Install]
WantedBy=multi-user.target
