# systemd/wlan0-ap.service
[Unit]
Description=Configure wlan0 for Access Point
Before=hostapd.service
After=NetworkManager.service network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes

# Configurar interfaz wlan0 (NetworkManager ya excluido en /etc/NetworkManager/conf.d/)
ExecStart=/sbin/ip addr flush dev wlan0
ExecStart=/sbin/ip addr add 192.168.50.1/24 dev wlan0
ExecStart=/sbin/ip link set wlan0 up

# NOTA: wlan0 solo alcanzará "state UP" cuando hostapd inicie el modo AP
# No verificamos aquí porque causaría fallo antes de que hostapd arranque

[Install]
WantedBy=multi-user.target
